# scripts/car-scraper.service
[Unit]
Description=Dealer Dev Ops API
After=network.target

[Service]
Type=simple
User=dealerdevops
Group=dealerdevops
WorkingDirectory=/opt/dealer-dev-ops
EnvironmentFile=/etc/dealer-dev-ops.env
ExecStart=/usr/bin/npm run start:pg
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
